# conf/systemd.service
[Unit]
Description=WG Manager API — __APP__
After=network.target

[Service]
Type=simple
# Root is required — wg-quick and wg commands need CAP_NET_ADMIN
User=root
WorkingDirectory=__INSTALL_DIR__/server
Environment=NODE_ENV=production
EnvironmentFile=__INSTALL_DIR__/server/.env
ExecStart=__NODE__ __INSTALL_DIR__/server/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=__APP__

LimitNOFILE=65535

# ── SYS-01: Hardening ─────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/wireguard __DATA_DIR__ /var/log /run
ProtectHome=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
# SystemCallFilter is intentionally omitted: Node.js requires a broad set of
# syscalls that cannot be reliably restricted without breaking the runtime.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
