# deploy/wg-man-api.service
[Unit]
Description=WG Manager API Server
After=network.target

[Service]
Type=simple
# Root is required — wg-quick and wg commands need CAP_NET_ADMIN
User=root
WorkingDirectory=/opt/wg-man/server
Environment=NODE_ENV=production
EnvironmentFile=/opt/wg-man/server/.env
ExecStart=/usr/bin/node /opt/wg-man/server/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-man

LimitNOFILE=65535

# SEC-13: systemd hardening — mirrors conf/systemd.service
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/wireguard /var/lib/wg-man /var/log /var/lib /run /var/run
ProtectHome=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
